# systemd/axp288-battery.service
[Unit]
Description=AXP288 battery level

[Service]
Type=simple
ExecStart=/opt/iwork8/scripts/battery.sh

[Install]
WantedBy=default.target
